# usb_gadget/gadget.service
[Unit]
Description=Audio Gadget

[Service]
ExecStart=/usr/local/bin/gadget.sh

[Install]
WantedBy=multi-user.target
